# services/syslog.service.conf
[program:syslog]
user=root
command=/usr/sbin/rsyslogd -d -n
